# perfetto-sys/deb/perfetto/etc/systemd/system/traced.service
[Unit]
Description=Traced Service
PartOf=perfetto.service

[Service]
ExecStart=/opt/perfetto/traced

[Install]
WantedBy=multi-user.target
